# templates/transparent_hugepage.conf
# Ubuntu upstart file at /etc/init/mongod_vm_settings.conf
#
# This file was generated by Ansible for {{ ansible_fqdn }}
# Do NOT modify this file by hand!

start on (starting mongod)
script
  echo "never" > /sys/kernel/mm/transparent_hugepage/enabled
  echo "never" > /sys/kernel/mm/transparent_hugepage/defrag
end script
